# scripts/launchd/com.andyherman.neural-bridge.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.andyherman.neural-bridge.discord-bot</string>

    <!-- caffeinate -s -i wraps the daemon to:
         (a) prevent process-scoped idle sleep (-i)
         (b) prevent system sleep on AC power (-s)
         The assertion is scoped to the daemon's lifetime — when the daemon
         exits, caffeinate exits, the assertion is released. This is on top
         of system-wide pmset settings as a belt-and-suspenders defense
         against macOS App Nap suspending background Python processes after
         long idle periods. -->
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/caffeinate</string>
        <string>-s</string>
        <string>-i</string>
        <string>/Users/andyherman/Development/neural-bridge/.venv/bin/python</string>
        <string>-m</string>
        <string>scripts.discord_bot.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andyherman/Development/neural-bridge</string>

    <!-- Start at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash, but don't tight-loop. -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- stdout / stderr. The daemon now configures its own RotatingFileHandler
         in main._configure_logging which writes to discord-bot.log (10MB cap,
         7 backups). launchd's StandardErrorPath captures anything that escapes
         the logger (boot-time crashes before _configure_logging runs, asyncio
         tracebacks, discord.py library output) — that file stays small in
         practice. StandardOutPath is unused (daemon doesn't print to stdout)
         so we don't redirect it. -->
    <key>StandardErrorPath</key>
    <string>/Users/andyherman/Library/Logs/neural-bridge/discord-bot.stderr.log</string>

    <!-- USER is required by keychain.py (the security CLI keys on it).
         PATH includes Homebrew bin so gh and claude resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>USER</key>
        <string>andyherman</string>
        <key>HOME</key>
        <string>/Users/andyherman</string>
        <key>PATH</key>
        <string>/Users/andyherman/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Process limits — sensible defaults. -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
