# scripts/launchd/com.andyherman.neural-bridge.discord-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.andyherman.neural-bridge.discord-bot</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/andyherman/Development/neural-bridge/.venv/bin/python</string>
        <string>-m</string>
        <string>scripts.discord_bot.main</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/andyherman/Development/neural-bridge</string>

    <!-- Start at login. -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Restart on crash, but don't tight-loop. -->
    <key>KeepAlive</key>
    <true/>
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <!-- stdout / stderr. main.py writes its log lines to stderr; flush.py and
         compile.py do not run inside this daemon, so the stderr file is the
         primary log surface. -->
    <key>StandardOutPath</key>
    <string>/Users/andyherman/Library/Logs/neural-bridge/discord-bot.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/andyherman/Library/Logs/neural-bridge/discord-bot.stderr.log</string>

    <!-- USER is required by keychain.py (the security CLI keys on it).
         PATH includes Homebrew bin so gh and claude resolve. -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>USER</key>
        <string>andyherman</string>
        <key>HOME</key>
        <string>/Users/andyherman</string>
        <key>PATH</key>
        <string>/Users/andyherman/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <!-- Process limits — sensible defaults. -->
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>5</integer>
</dict>
</plist>
